# scripts/hestia-pir.service
[Unit]
Description=Hestia PIR motion sensor bridge
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple

# !!! ANPASSEN !!!
# User= muss mit deinem echten Username im /etc/passwd übereinstimmen,
# sonst exit 217/USER. Default 'pi' aus dem alten RPi-OS-Imager — viele
# Systeme verwenden heute einen eigenen Namen.
User=pi
WorkingDirectory=%h/hestia

# Selbst-heilendes execute-Bit (siehe hestia-kiosk.service für die ausführliche
# Erklärung). `-` macht den Aufruf failsoft.
ExecStartPre=-/bin/chmod +x %h/hestia/scripts/pir-motion.py

ExecStart=/usr/bin/python3 %h/hestia/scripts/pir-motion.py
Environment=HESTIA_PIR_GPIO=17
Environment=HESTIA_MOTION_URL=http://localhost:3001/api/internal/motion
# Wenn MOTION_SECRET in der Backend-.env gesetzt ist, hier identisch eintragen:
# Environment=HESTIA_MOTION_SECRET=dein-secret
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
