# scripts/hestia-kiosk.service
[Unit]
Description=Hestia Wall Display (Chromium kiosk via cage)
After=network-online.target hestia-pir.service
Wants=network-online.target

[Service]
Type=simple

# !!! ANPASSEN !!!
# User= MUSS der Benutzer sein, den du beim RPi-OS-Imager angegeben hast.
# Im offiziellen RPi-Imager-Default ist es 'pi', viele wählen einen anderen
# Namen (z.B. via "OS-Anpassung → Benutzername"). Stimmt der Wert hier nicht
# mit einem echten /etc/passwd-Eintrag überein, scheitert systemd beim Start
# mit `status=217/USER`. Entweder direkt hier editieren ODER per Override:
#   sudo systemctl edit hestia-kiosk
#   [Service]
#   User=<dein-username>
User=pi

PAMName=login
TTYPath=/dev/tty1
StandardInput=tty
StandardOutput=journal
StandardError=journal

# Selbst-heilendes execute-Bit: Wenn das Repo via SCP/rsync von Windows
# auf den RPi gekommen ist, gehen file modes häufig verloren — die
# Skripte landen mit 0644 (no-exec), cage stirbt dann mit "Failed to
# spawn client: Permission denied". Mit `chmod +x` als ExecStartPre
# bügeln wir das vor jedem Start aus. `-` macht den Aufruf failsoft —
# wenn die Datei fehlt, kein Service-Abort.
ExecStartPre=-/bin/chmod +x %h/hestia/scripts/hestia-kiosk.sh %h/hestia/scripts/hestia-kiosk-login.sh

# Cage = minimaler Wayland-Kiosk-Compositor, startet genau ein Programm fullscreen.
# Bildschirmschoner deaktivieren wir bewusst nicht hier — useWallScreensaver
# im Frontend übernimmt das (dimmt zu schwarz, weckt auf bei PIR/Touch).
ExecStart=/usr/bin/cage -s -- %h/hestia/scripts/hestia-kiosk.sh

Environment=HESTIA_WALL_URL=https://hestia.local/wall

# Optional: Auto-Login für den Kiosk. Wenn diese Datei existiert (chmod 600!),
# loggt der Kiosk beim Start automatisch ein. Format der Datei:
#   HESTIA_KIOSK_EMAIL=wand@hestia.local
#   HESTIA_KIOSK_PASSWORD=dein-passwort
# Ohne diese Datei zeigt der Kiosk das Login-Formular und ein User muss
# einmalig manuell einloggen (Token bleibt dann 30 Tage gültig).
EnvironmentFile=-%h/.hestia-kiosk.env

# /run/user/<uid> für Wayland-Sockets. Voraussetzung:
# `sudo loginctl enable-linger <user>` einmalig ausführen, damit das
# Verzeichnis ohne interaktiven Login existiert. Die hartkodierte 1000
# stimmt mit dem default RPi-OS-User überein; falls dein User eine andere
# UID hat, hier anpassen (`id -u <user>`).
Environment=XDG_RUNTIME_DIR=/run/user/1000

# Auf RPi mit DSI-Display (Touch Display) sitzt der Display-fähige DRM-
# Node oft auf card1, nicht card0. Wenn dein Setup HDMI nutzt, ist es
# umgekehrt — dann diese Zeile auf card0 setzen oder ganz löschen.
# Prüfen mit:
#   for f in /sys/class/drm/card*-*/status; do echo "$f: $(cat $f)"; done
Environment=WLR_DRM_DEVICES=/dev/dri/card1

Restart=always
RestartSec=3

# Wenn der RPi nichts anderes auf tty1 zeigen soll:
# Stelle sicher, dass getty@tty1 deaktiviert ist (siehe DEPLOYMENT-RPI.md).

[Install]
WantedBy=multi-user.target
